# aur-packages/meowrch-settings/usr/lib/systemd/system/sddm-plymouth.service.d/00-meowrch-timeout.conf
[Service]
# Increase timeout for sddm-plymouth to prevent restart loops
# This is especially important when system optimizations are applied
TimeoutStartSec=120s
TimeoutStopSec=30s

# Prevent rapid restarts that can cause login loops
StartLimitInterval=300s
StartLimitBurst=3
