# aur-packages/meowrch-settings/usr/lib/systemd/system/sddm-plymouth.service.d/00-meowrch-timeout.conf
[Service]
# Increase the timeout for sddm-plymouth to prevent restart cycles
TimeoutStartSec=120s
TimeoutStopSec=30s

# Prevent rapid restarts that can cause login loops
StartLimitInterval=300s
StartLimitBurst=3
