# home/.config/systemd/user/firefox-theme-update.service
[Unit]
Description=Firefox GNOME Theme Auto-Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-firefox-gnome-theme.sh
