# home/.config/systemd/user/firefox-theme-update.timer
[Unit]
Description=Firefox GNOME Theme Auto-Update Timer
Requires=firefox-theme-update.service

[Timer]
OnBootSec=2min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
